# tools/demo-backup/ofbizDemo.crontab
0 3 * * * /home/ofbizDemo/check-svn-update.sh > /home/ofbizDemo/cronlog-svn-update.log 2>&1