# debian/ofbiz-specialpurpose.postinst
#!/bin/sh
set -e

. /usr/share/debconf/confmodule
ofbiz_init_invoke() {
	if [ -x "/etc/init.d/ofbiz" ]; then
		if [ -x "`which invoke-rc.d 2>/dev/null`" ]; then
			invoke-rc.d ofbiz $1 || exit $?
		else
			/etc/init.d/ofbiz $1 || exit $?
		fi
	fi
}
case "$1" in
	(configure|reconfigure)
		ofbiz_init_invoke stop
		dpkg-trigger --no-await ofbiz-install
		;;
esac
#DEBHELPER#
db_stop
